# daemon/armor-autostart.service.in
[Unit]
Description="Autostart/restart armor"
[Service]
WorkingDirectory=${CMAKE_CURRENT_BINARY_DIR}
ExecStart=/usr/bin/xvfb-run ${CMAKE_CURRENT_BINARY_DIR}/attack_mind
Restart=always
[Install]
WantedBy=multi-user.target
